# luci-app-xray-client/files/etc/init.d/xrayclient
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

XRAY_DIR="/usr/share/xrayclient"
CONF_DIR="/etc/xrayclient/conf.d"
LOG_FILE="/var/log/xrayclient.log"
UCI_CONF="xrayclient"

start_service() {
    # 确保日志目录存在
    mkdir -p /var/log 2>/dev/null
    
    # 使用 { ... } 将整个启动过程打包，统一输出到日志文件
    {
        echo "===== $(date '+%Y-%m-%d %H:%M:%S') 开始启动 Xray client ====="

        # 从 xray-core 的 UCI 配置读取 datadir (geoip/geosite 所在目录)
        ASSET_DIR=$(uci -q get xray.config.datadir)
        [ -z "$ASSET_DIR" ] && ASSET_DIR="/usr/share/xray"

        # 在本目录下创建符号链接供前端读取 (前端 ACL 无法访问动态路径)
        ln -sf "$ASSET_DIR/geoip.dat" "$XRAY_DIR/geoip.dat"
        ln -sf "$ASSET_DIR/geosite.dat" "$XRAY_DIR/geosite.dat"

        # 同步定时更新任务 (IP类 和 域名类 分别调度)
        mkdir -p /etc/crontabs
        touch /etc/crontabs/root
        sed -i '/# xrayclient-update-ip\|# xrayclient-update-dat/d' /etc/crontabs/root 2>/dev/null

        # --- IP 类更新调度 (cn_v4 + cn_v6 + geoip) ---
        IP_INTERVAL=$(uci -q get ${UCI_CONF}.main.ip_update_interval)
        [ -z "$IP_INTERVAL" ] && IP_INTERVAL="weekly"
        IP_HOUR=$(uci -q get ${UCI_CONF}.main.ip_update_hour)
        [ -z "$IP_HOUR" ] && IP_HOUR="4"
        IP_DOW=$(uci -q get ${UCI_CONF}.main.ip_update_dow)
        [ -z "$IP_DOW" ] && IP_DOW="5"

        case "$IP_INTERVAL" in
            daily)
                echo "0 ${IP_HOUR} * * * sh /usr/share/xrayclient/update_data.sh ip # xrayclient-update-ip" >> /etc/crontabs/root
                ;;
            every3d)
                echo "0 ${IP_HOUR} */3 * * sh /usr/share/xrayclient/update_data.sh ip # xrayclient-update-ip" >> /etc/crontabs/root
                ;;
            weekly)
                echo "0 ${IP_HOUR} * * ${IP_DOW} sh /usr/share/xrayclient/update_data.sh ip # xrayclient-update-ip" >> /etc/crontabs/root
                ;;
            never)
                echo "IP 类数据定时更新已禁用" >> "$LOG_FILE"
                ;;
        esac

        # --- 域名类更新调度 (geosite) ---
        DAT_INTERVAL=$(uci -q get ${UCI_CONF}.main.dat_update_interval)
        [ -z "$DAT_INTERVAL" ] && DAT_INTERVAL="every3d"
        DAT_HOUR=$(uci -q get ${UCI_CONF}.main.dat_update_hour)
        [ -z "$DAT_HOUR" ] && DAT_HOUR="4"
        DAT_DOW=$(uci -q get ${UCI_CONF}.main.dat_update_dow)
        [ -z "$DAT_DOW" ] && DAT_DOW="1"

        case "$DAT_INTERVAL" in
            daily)
                echo "0 ${DAT_HOUR} * * * sh /usr/share/xrayclient/update_data.sh dat # xrayclient-update-dat" >> /etc/crontabs/root
                ;;
            every3d)
                echo "0 ${DAT_HOUR} */3 * * sh /usr/share/xrayclient/update_data.sh dat # xrayclient-update-dat" >> /etc/crontabs/root
                ;;
            weekly)
                echo "0 ${DAT_HOUR} * * ${DAT_DOW} sh /usr/share/xrayclient/update_data.sh dat # xrayclient-update-dat" >> /etc/crontabs/root
                ;;
            never)
                echo "域名类数据定时更新已禁用" >> "$LOG_FILE"
                ;;
        esac

        /etc/init.d/cron enable 2>/dev/null
        /etc/init.d/cron restart 2>/dev/null

        # 检查用户是否已初始化，未初始化则先添加用户
        USR_SETUP=$(uci -q get ${UCI_CONF}.main.usr_setup)
        if [ "$USR_SETUP" != "1" ]; then
            echo "首次启动: 正在添加 Xray 运行用户..."
            sh "$XRAY_DIR/add_usr.sh"
            if [ $? -ne 0 ]; then
                echo "错误: add_usr.sh 执行失败，无法启动服务。"
                return 1
            fi
        fi

        # 检查是否配置了活动节点
        ACTIVE_NODE=$(uci -q get ${UCI_CONF}.main.active_node)
        if [ -z "$ACTIVE_NODE" ]; then
            echo "未配置活动节点 (active_node)，跳过服务启动。"
            echo "请在 LuCI 页面选择代理节点后再启动服务。"
            return 0
        fi

        XRAY_BIN="/usr/bin/xray"
        if [ ! -f "$XRAY_BIN" ]; then
            echo "错误: 找不到 Xray 二进制文件 $XRAY_BIN"
            return 1
        fi

        # 从 UCI 读取运行用户配置
        XRAY_USR=$(uci -q get ${UCI_CONF}.main.xray_usr_name)
        XRAY_GID=$(uci -q get ${UCI_CONF}.main.xray_usr_gid)
        
        if [ -z "$XRAY_USR" ] || [ -z "$XRAY_GID" ]; then
            echo "xrayclient: UCI 配置异常，缺少 xray_usr_name 或 xray_usr_gid。"
            return 1
        fi

        # 防御性清理旧规则
        [ -f "$XRAY_DIR/remove_nft.sh" ] && sh "$XRAY_DIR/remove_nft.sh" > /dev/null 2>&1
        [ -f "$XRAY_DIR/remove_route.sh" ] && sh "$XRAY_DIR/remove_route.sh" > /dev/null 2>&1

        # 生成 Xray 配置
        sh "$XRAY_DIR/gen_xray_config.sh"
        if [ $? -ne 0 ]; then
            echo "错误: gen_xray_config.sh 执行失败"
            return 1
        fi

        # 启动 xray 主进程
        procd_open_instance
        procd_set_param command "$XRAY_BIN" -confdir "$CONF_DIR"
        procd_set_param user "$XRAY_USR"
        procd_set_param group "$XRAY_GID"
        procd_set_param env XRAY_LOCATION_ASSET="$ASSET_DIR"
        procd_set_param limits nofile=1000000 nproc=500
        procd_set_param respawn
        procd_set_param stdout 1
        procd_set_param stderr 1
        procd_close_instance

        # 添加路由策略
        sh "$XRAY_DIR/add_route.sh"
        if [ $? -ne 0 ]; then
            echo "错误: add_route.sh 执行失败，执行回滚..."
            stop
            return 1
        fi

        # 添加 nftables 规则
        sh "$XRAY_DIR/add_nft.sh"
        if [ $? -ne 0 ]; then
            echo "错误: add_nft.sh 执行失败，执行回滚..."
            stop
            return 1
        fi

        echo "===== $(date '+%Y-%m-%d %H:%M:%S') Xray client 启动成功！ ====="
    } >> "$LOG_FILE" 2>&1
}

stop_service() {
    mkdir -p /var/log 2>/dev/null
    
    # 将停止过程也输出到日志
    {
        echo "===== $(date '+%Y-%m-%d %H:%M:%S') 开始停止 Xray client ====="
        
        # 依次执行 remove
        [ -f "$XRAY_DIR/remove_nft.sh" ] && sh "$XRAY_DIR/remove_nft.sh"
        [ -f "$XRAY_DIR/remove_route.sh" ] && sh "$XRAY_DIR/remove_route.sh"
        
        echo "===== $(date '+%Y-%m-%d %H:%M:%S') Xray client 已停止。 ====="
    } >> "$LOG_FILE" 2>&1
}

reload_service() {
    stop
    start
}
